# api/stress.cron
0 0 * * * python /usr/src/app/stres_clac_batch.py
